# ansible/roles/vault_agent/templates/vault-agent.service.j2
# /etc/systemd/system/vault-agent-{{ vault_agent_name }}.service
# Managed by ansible/roles/vault_agent.

[Unit]
Description=Vault Agent for {{ vault_agent_name }}
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
Group=root
{# Recreate the consumer's volatile (/run) dirs on boot — tmpfs is wiped.
   Persistent dirs (e.g. /etc/*) are handled by the role's file task. #}
{% for d in vault_agent_extra_dirs if d.path.startswith('/run/') %}
ExecStartPre=/bin/mkdir -p {{ d.path }}
ExecStartPre=/bin/chgrp {{ d.group | default('root') }} {{ d.path }}
ExecStartPre=/bin/chmod {{ d.mode | default('0750') }} {{ d.path }}
{% endfor %}
ExecStart=/usr/bin/vault agent -config={{ vault_agent_config_dir }}/{{ vault_agent_name }}.hcl
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=5
RuntimeDirectory=vault-agent
RuntimeDirectoryMode=0750
NoNewPrivileges=yes
PrivateTmp=yes
ProtectHome=yes
ProtectKernelModules=yes
ProtectKernelTunables=yes

[Install]
WantedBy=multi-user.target
