# ansible/generated/extmon/extmon-diag-agent.service
[Unit]
Description=AS215932 external diagnostic probe agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=prometheus
Group=prometheus
EnvironmentFile=/etc/extmon/diag-agent.env
ExecStart=/usr/local/bin/extmon-diag-agent.py
Restart=always
RestartSec=10s
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true

[Install]
WantedBy=multi-user.target
