# configs/rtr/dnat-vrf-leak.service
# /etc/systemd/system/dnat-vrf-leak.service
# Route leak: 10.0.2.0/24 in main table so DNAT'd IPv4 traffic
# (arriving on enX4, default VRF) can reach infra VMs on enX2 (overlay VRF).
[Unit]
Description=DNAT VRF route leak (IPv4 infra subnet)
After=systemd-networkd.service
Before=nftables.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/sbin/ip route replace 10.0.2.0/24 dev enX2
ExecStop=/usr/sbin/ip route del 10.0.2.0/24 dev enX2

[Install]
WantedBy=multi-user.target
