# ansible/roles/knot/templates/knot-online-backup.service.j2
[Unit]
Description=Online backup of Knot zones, catalogs, and DNSSEC keys
After={{ knot_service }}.service
Requires={{ knot_service }}.service

[Service]
Type=oneshot
User=knot
Group=knot
ExecStart=/usr/local/sbin/knot-online-backup
Nice=10
IOSchedulingClass=idle
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths={{ knot_backup_dir }}
