# ansible/roles/soc_agent/templates/soc-handoffs.timer.j2
[Unit]
Description=Poll central coordinator for inbound SOC handoffs

[Timer]
OnBootSec=30
OnUnitActiveSec={{ soc_handoff_on_unit_active_sec }}s
AccuracySec=1s
Persistent=false

[Install]
WantedBy=timers.target
