# ansible/generated/cr1-de1/softflowd_ixp_fra
#!/bin/sh

# PROVIDE: softflowd_ixp_fra
# REQUIRE: NETWORKING
# KEYWORD: shutdown

. /etc/rc.subr

name="softflowd_ixp_fra"
rcvar="${name}_enable"
pidfile="/var/run/${name}.pid"
ctlfile="/var/run/${name}.ctl"
command="/usr/sbin/daemon"
softflowd_binary="/usr/local/sbin/softflowd"
softflowd_args="-d -N -i vtnet3 -n [2a0c:b641:b50:2::110]:2055 -v 9 -s 1000 -m 65536 -T full -e 2a0c:b641:b50::b"
command_args="-f -p ${pidfile} ${softflowd_binary} ${softflowd_args} -c ${ctlfile}"
stop_cmd="${name}_stop"

load_rc_config "$name"
: ${softflowd_ixp_fra_enable:="NO"}

softflowd_ixp_fra_stop()
{
    if [ -S "${ctlfile}" ]; then
        /usr/local/bin/softflowctl -c "${ctlfile}" exit >/dev/null 2>&1 || true
    fi
    if [ -f "${pidfile}" ]; then
        kill "$(cat "${pidfile}")" >/dev/null 2>&1 || true
        rm -f "${pidfile}"
    fi
}

run_rc_command "$1"
