# ansible/roles/agent_mail/templates/agent-mail-backup.service.j2
[Unit]
Description=Quiesced Stalwart Agent Mail backup
Requires=docker.service
RequiresMountsFor={{ agent_mail_backup_dir }}
After=docker.service

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/agent-mail-backup
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7
PrivateTmp=true
ProtectHome=true
