# configs/loop/hyrule-engineering-loop.timer
# /etc/systemd/system/hyrule-engineering-loop.timer
# Deploy to: the same operator machine / `loop` VM as the service.
#
# Hourly cadence; the daemon's own per-run and per-day budgets (and the run
# lock) bound the work, so the timer only needs to fire often enough to keep
# the Icinga freshness check (check_interval = 90m) green. Keep this cadence
# and the icinga engineering-loop.conf freshness window in step.

[Unit]
Description=Schedule the AS215932 Engineering Loop operations lane

[Timer]
OnCalendar=hourly
# Avoid stampeding other top-of-hour jobs.
RandomizedDelaySec=300
Persistent=true
Unit=hyrule-engineering-loop.service

[Install]
WantedBy=timers.target
