# configs/loop/hyrule-engineering-loop.service
# /etc/systemd/system/hyrule-engineering-loop.service
# Deploy to: an operator machine or a dedicated `loop` VM — NEVER the
# privileged `ci` runner. The backend executes generated code; the daemon
# refuses to run when GITHUB_ACTIONS is set, but do not schedule it there
# regardless. A dedicated `loop` VM goes through the standard
# docs/network-flows.md + firewall + monitoring onboarding three-steps.
#
# Oneshot, timer-driven: one cycle picks one loop:approved issue, runs the
# full graph, and ends with a draft PR or a journaled failure, then exits.

[Unit]
Description=AS215932 Engineering Loop operations lane (one budgeted cycle)
After=network-online.target
Wants=network-online.target
# Belt and braces: do not even attempt a cycle on a CI runner.
ConditionEnvironment=!GITHUB_ACTIONS

[Service]
Type=oneshot
User=loop
Group=loop
WorkingDirectory=/opt/engineering-loop
EnvironmentFile=/opt/engineering-loop/.env
# .env provides: HYRULE_DISCORD_WEBHOOK, HYRULE_ICINGA_URL/USER/PASSWORD,
# the loop's GH token for `gh`, and any model/provider keys. It must NOT
# contain production Vault tokens or fleet SSH material — the backend env
# is scrubbed, but keep the surface minimal.
ExecStart=/opt/engineering-loop/.venv/bin/hyrule-engineering-loop daemon --once \
    --workspace-root /opt/engineering-loop/workspace \
    --output-root /var/lib/engineering-loop/runs \
    --state-dir-path /var/lib/engineering-loop/state \
    --memory-dir /opt/engineering-loop/workspace/hyrule-infra/memory
# A run never blocks the next timer fire indefinitely; the daemon enforces
# its own per-run wall-clock budget, this is the hard backstop.
TimeoutStartSec=3600
StandardOutput=journal
StandardError=journal
SyslogIdentifier=engineering-loop

# Security hardening
NoNewPrivileges=yes
PrivateTmp=yes
ProtectHome=yes
ProtectKernelModules=yes
ProtectKernelTunables=yes
ProtectSystem=strict
ReadWritePaths=/opt/engineering-loop /var/lib/engineering-loop

[Install]
WantedBy=multi-user.target
