# ansible/roles/hyrule_network_proxy/templates/hyrule-network-proxy.service.j2
# /etc/systemd/system/hyrule-network-proxy.service
# Managed by Ansible.

[Unit]
Description=Hyrule Network Proxy
After=network-online.target tor.service i2pd.service yggdrasil.service
Wants=network-online.target tor.service i2pd.service yggdrasil.service

[Service]
Type=exec
User={{ hyrule_network_proxy_user }}
Group={{ hyrule_network_proxy_group }}
EnvironmentFile={{ hyrule_network_proxy_env_file }}
ExecStart={{ hyrule_network_proxy_bin_path }}
Restart=always
RestartSec=5
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths={{ hyrule_network_proxy_state_dir }}
StandardOutput=journal
StandardError=journal
SyslogIdentifier=hyrule-network-proxy

[Install]
WantedBy=multi-user.target
