# ansible/generated/agentmail/agent-mail-backup.service
[Unit]
Description=Quiesced Stalwart Agent Mail backup
Requires=docker.service
After=docker.service

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/agent-mail-backup
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7
PrivateTmp=true
ProtectHome=true
