# configs/hyrule-cloud-worker.service
# /etc/systemd/system/hyrule-cloud-worker.service
# Exactly one worker consumes durable domain/native-payment jobs. API workers
# never run schedulers in-process.

[Unit]
Description=Hyrule Cloud durable background worker
After=network-online.target postgresql.service hyrule-cloud.service monero-wallet-rpc.service
Requires=postgresql.service
Wants=hyrule-cloud.service monero-wallet-rpc.service

[Service]
Type=exec
User=hyrule
Group=hyrule
WorkingDirectory=/opt/hyrule-cloud
EnvironmentFile=/opt/hyrule-cloud/.env
Environment=PYTHONPATH=/opt/hyrule-cloud
ExecStart=/opt/hyrule-cloud/.venv/bin/hyrule-cloud-worker
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=hyrule-cloud-worker

ProtectSystem=strict
ReadWritePaths=/opt/hyrule-cloud /var/lib/hyrule-cloud
NoNewPrivileges=yes
PrivateTmp=yes
ProtectHome=yes
ProtectKernelModules=yes
ProtectKernelTunables=yes

[Install]
WantedBy=multi-user.target
