# configs/noc-agent.service
# /etc/systemd/system/noc-agent.service
# Deploy to: noc VM

[Unit]
Description=AS215932 NOC Agent (Alertmanager + Icinga webhooks → PydanticAI investigation → Discord)
After=network.target postgresql.service hyrule-mcp.service xo-mcp.service
Wants=postgresql.service hyrule-mcp.service xo-mcp.service

[Service]
Type=exec
User=noc-agent
Group=noc-agent
WorkingDirectory=/opt/noc-agent
EnvironmentFile=/opt/noc-agent/.env
Environment=PYTHONPATH=/opt/noc-agent
ExecStartPre=/bin/grep -Eq "^NOC_APPROVAL_SIGNING_SECRET=.{32,}$" /opt/noc-agent/.env
ExecStartPre=/bin/grep -Eq "^HYRULE_MCP_ACTION_SIGNING_SECRET=.{32,}$" /opt/noc-agent/.env
ExecStart=/opt/noc-agent/.venv/bin/uvicorn app.main:app \
    --host :: \
    --port 8000 \
    --workers 2 \
    --log-level info \
    --access-log
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=noc-agent

# Security hardening
ProtectSystem=strict
ReadWritePaths=/opt/noc-agent /var/lib/noc-agent
ReadOnlyPaths=/etc/hyrule-mcp /var/lib/hyrule-mcp
NoNewPrivileges=yes
PrivateTmp=yes
ProtectHome=yes
ProtectKernelModules=yes
ProtectKernelTunables=yes

[Install]
WantedBy=multi-user.target
