# ansible/roles/github_runner/templates/github-runner.service.j2
# Managed by ansible/roles/github_runner
[Unit]
Description=GitHub Actions runner ({{ github_runner_name }})
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User={{ github_runner_user }}
Group={{ github_runner_group }}
WorkingDirectory={{ github_runner_install_dir }}
ExecStart={{ github_runner_install_dir }}/run.sh
# run.sh exits 0 when the runner self-updates or is reconfigured; on-failure
# would leave the service dead after that. always — matches GitHub's own
# generated unit — so the runner comes back regardless of exit code.
Restart=always
RestartSec=5s
# mixed (NOT GitHub's default `process`): SIGTERM to run.sh, then SIGKILL to the
# whole control group at TimeoutStopSec. With `process`, a restart mid-job (disk
# migration, vault-agent env re-render, manual restart) orphaned Runner.Listener:
# it kept holding the GitHub session — so the next start hit SessionConflictException
# and crash-looped (NRestarts 245+) — and ran jobs inside the torn-down PrivateTmp
# namespace, where `mktemp /tmp/...` failed ("No such file or directory"). `mixed`
# reaps the listener together with the unit, so the session releases and /tmp is
# never stale. Trade-off: a runner self-update during an in-flight job kills that
# job (it re-queues) — rare and cheap next to the orphan loop.
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=5min
Environment=HOME={{ github_runner_home }}
Environment=TMPDIR={{ github_runner_tmp_dir }}

# Hardening
NoNewPrivileges=true
ProtectSystem=full
ProtectHome=read-only
ReadWritePaths={{ github_runner_install_dir }} {{ github_runner_home }}
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

[Install]
WantedBy=multi-user.target
