# ansible/roles/hyrule_mcp/templates/bgp-router-snapshot.service.j2
[Unit]
Description=Collect AS215932 router BGP table snapshots
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User={{ hyrule_mcp_app_user }}
Group={{ hyrule_mcp_app_group }}
WorkingDirectory={{ hyrule_mcp_install_dir }}
Environment=HYRULE_MCP_HOSTS={{ hyrule_mcp_etc_dir }}/hosts.yml
EnvironmentFile={{ hyrule_mcp_etc_dir }}/bgp-router-snapshot.env
ExecStart={{ hyrule_mcp_install_dir }}/.venv/bin/hyrule-mcp-bgp-snapshot
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths={{ hyrule_mcp_state_dir }}
